# autostart/plank-lxqt.desktop
[Desktop Entry]
Name=Plank-LXQt
Categories=Utility;
Type=Application
Exec=plank-lxqt
Icon=plank
Terminal=false
NoDisplay=false
OnlyShowIn=LXQt
