# autostart/plank-lxde.desktop
[Desktop Entry]
Name=Plank-LXDE
Categories=Utility;
Type=Application
Exec=plank-lxde
Icon=plank
Terminal=false
NoDisplay=false
OnlyShowIn=LXDE

