# autostart/plank-razor.desktop
[Desktop Entry]
Name=Plank-Razor
Categories=Utility;
Type=Application
Exec=plank-razor
Icon=plank
Terminal=false
NoDisplay=false
OnlyShowIn=Razor

