# autostart/plank-xfce.desktop
[Desktop Entry]
Name=Plank-XFCE
Categories=Utility;
Type=Application
Exec=plank-xfce
Icon=plank
Terminal=false
NoDisplay=false
OnlyShowIn=XFCE
