# autostart/plank-gnome.desktop
[Desktop Entry]
Name=Plank-GNOME
Categories=Utility;
Type=Application
Exec=plank-gnome
Icon=plank
Terminal=false
NoDisplay=false
OnlyShowIn=GNOME
